# docker/supervisor/pingos.conf
[program:nginx]
command=bash pingos-entrypoint.sh
directory=/usr/local/pingos/
autostart=true
autorestart=true
redirect_stderr=true
priority=10
stdout_logfile=/var/log/supervisor/pingos.log